# blog-store-messages.service.in
[Unit]
Description=Tell blogd about writable boot log
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target
ConditionPathExists=!/etc/initrd-release

[Service]
Type=oneshot
TimeoutSec=20
RemainAfterExit=yes
ExecStart=-/sbin/blogctl ready
