# blog-store-messages.service.in
[Unit]
Description=Tell blogd about writable boot log
DefaultDependencies=no
RequiresMountsFor=/var/log/boot.log
Conflicts=shutdown.target
After=systemd-remount-fs.service systemd-tmpfiles-setup.service auditd.service
After=local-fs.target
Before=sysinit.target
ConditionPathExists=!/etc/initrd-release
ConditionPathIsReadWrite=/var/log

[Service]
Type=oneshot
TimeoutSec=20
RemainAfterExit=yes
ExecStart=-/sbin/blogctl ready
